# pyess-addon/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the pyess-addon service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

bashio::log.info "declaring variables"

# Declare variables
declare ess_host sensors ess_password mqtt_host mqtt_user mqtt_password

## Get the 'message' key from the user config options.

bashio::log.info "reading config"

ess_host=$(bashio::config 'ess_host')
sensors=$(bashio::config 'ess_sensors')
ess_password=$(bashio::config 'ess_password')
mqtt_host=$(bashio::services mqtt "host")
mqtt_user=$(bashio::services mqtt "username")
mqtt_password=$(bashio::services mqtt "password")

bashio::log.info "done reading config"


## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Starting pyss with ess ${ess_host}, logging to mqtt server ${mqtt_host}"

## Run your program
/usr/bin/my_program
/usr/bin/essmqtt --ess_password ${ess_password} --mqtt_server ${mqtt_host} --mqtt_user=${mqtt_user} --mqtt_password=${mqtt_password} --interval_seconds=10 --ess_host=${ess_host} --hass_autoconfig_sensors=${sensors}
